# crontab
# Еженедельный отчет по понедельникам в 10:00
0 10 * * 1 cd /Users/alekenov/CascadeProjects/seobot && python3 src/scripts/send_weekly_report.py >> logs/weekly_report.log 2>&1

# Ежедневный сбор статистики в 9:00
0 9 * * * cd /Users/alekenov/CascadeProjects/seobot && python3 src/scripts/collect_daily_stats.py >> logs/daily_stats.log 2>&1
